# supervisor.ini
[supervisord]
logfile = /var/log/supervisord.log
logfile_backups = 0

[program:rsyslogd]
command = /usr/local/bin/daemonize.sh /var/run/rsyslog.pid /usr/sbin/rsyslogd -n

[program:postfix]
command = /usr/local/bin/daemonize.sh /var/spool/postfix/pid/master.pid /usr/sbin/postfix start
startsecs = 3
